# scripts/com.github.pajaksmart.autopush.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.github.pajaksmart.autopush</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>python3</string>
    <string>/Users/awantjhai/Library/CloudStorage/GoogleDrive-halo@pajaksmart.com/My Drive/AI-Pajaksmart/websites/pajaksmart.xyz/scripts/auto_push.py</string>
    <string>--path</string>
    <string>/Users/awantjhai/Library/CloudStorage/GoogleDrive-halo@pajaksmart.com/My Drive/AI-Pajaksmart/websites/pajaksmart.xyz</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/autopush.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/autopush.err.log</string>
</dict>
</plist>
